# quickshell/scripts/aikira/aikira.service
[Unit]
Description=Aikira AI Chat Backend
After=network.target postgresql.service

[Service]
Type=simple
# ── Edit these three paths ──────────────────────────────────────────────────
WorkingDirectory=%h/aikira
ExecStart=%h/aikira/.venv/bin/python run.py
EnvironmentFile=%h/aikira/.env
# ────────────────────────────────────────────────────────────────────────────

Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=aikira

[Install]
WantedBy=default.target
